# deploy/stanley/templates/stanley-api.service.j2
[Unit]
Description=Stanley Financial Platform API
Documentation=https://github.com/your-org/stanley
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=exec
User={{ stanley_user }}
Group={{ stanley_user }}
WorkingDirectory={{ stanley_home }}
Environment=PATH={{ stanley_home }}/venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile={{ stanley_home }}/.env
ExecStart={{ stanley_home }}/venv/bin/python -m stanley.api.main
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths={{ stanley_data_dir }} {{ stanley_logs_dir }}
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=1G
TasksMax=4096

[Install]
WantedBy=multi-user.target